# roles/instance_driver/templates/audiocloud_instance_driver.service.j2
[Unit]
Description=AudioCloud Instance Driver - {{ instance_driver_version }}
Documentation=https://audiocloud.io
After=network.target

[Service]
Environment=NATS_URL={{ instance_driver_nats_url }}
Type=simple
User={{ user }}
ExecStart={{ instance_driver_path }}/instance_driver driver.yaml
Restart=on-failure

[Install]
WantedBy=multi-user.target
